# onedrive-player.service
[Unit]
Description=OneDrive Video Player
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/onedrive-player
ExecStart=/usr/bin/python3 /opt/onedrive-player/player_app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=AZURE_CONFIG=/opt/onedrive-player/azure_config.json
Environment=ONEDRIVE_FOLDER=SOOP_VOD

[Install]
WantedBy=multi-user.target
